# vault.service
[Unit]
Description=Vault Server

[Service]
ExecStart=vault server -config=/tmp/vault_config.hcl